# launchagents/com.vadims.vadimgest-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vadims.vadimgest-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>-m</string>
        <string>vadimgest</string>
        <string>serve</string>
        <string>--no-open</string>
        <string>--port</string>
        <string>8484</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__/vadimgest</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/vadimgest-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/vadimgest-dashboard.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:__PYTHON_DIR__:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>GOG_KEYRING_PASSWORD</key>
        <string>klava-default-keyring</string>
    </dict>
</dict>
</plist>
