# launchagents/com.vadims.cron-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vadims.cron-scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>__REPO_DIR__/gateway/cron-scheduler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__/gateway</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cron-scheduler.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cron-scheduler.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:__PYTHON_DIR__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>CLAUDE_CONFIG_DIR</key>
        <string>__REPO_DIR__/.claude</string>
        <key>GOG_KEYRING_PASSWORD</key>
        <string>klava-default-keyring</string>
    </dict>
</dict>
</plist>
